# examples/dontwi.service
[Unit]
Description=Dontwi

[Service]
Type=simple
ExecStart=/opt/dontwi/dontwi.sh
WorkingDirectory=/opt/dontwi

[Install]
WantedBy=default.target
